# src/systemd/ubuntu-transition-tracker.timer
[Unit]
Description=Ubuntu Transition Tracker

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
